# build/homebridge/root/etc/services.d/homebridge/run
#!/usr/bin/with-contenv bash

echo "Starting Homebridge."
exec \
	s6-setuidgid abc /bin/bash -c \
	'/usr/bin/homebridge -D -P $plugin_folder -U /homebridge-config'